# ls-hubd_public.service
[Unit]
Description=Luna Service - private bus
Requires=PmLogDaemon.service

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/sbin/ls-hubd --public --conf /etc/luna-service2/ls-public.conf --pmloglib

[Install]
WantedBy=webos-essential.target
