# LunaSysMgr.service
[Unit]
Description=Luna System Manager

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/sbin/LunaSysMgr

[Install]
WantedBy=webos-base.target
