# native-applications.slice
[Unit]
Description=Slice for native applications
Wants=-.slice
After=-.slice
