# web-applications.slice
[Unit]
Description=Slice for web based applications
Wants=-.slice
After=-.slice
