# webos-base.target
[Unit]
Description=Starts up the basic webOS system daemons

[Install]
WantedBy=multi-user.target
