# webos-essential.target
[Unit]
Description=Starts up the essential system daemons for the webOS system
